# takoserverd.service
[Unit]
Description=takoserver

[Service]
ExecStart=/opt/takoserver/bin/takoserver
EnvironmentFile=/etc/default/takoserver
Restart=always
Type=simple
User=takoyaki
Group=takoyaki

[Install]
WantedBy=multi-user.target
